# 90-fstab-automount.rules
ACTION=="add|change|remove",SUBSYSTEM=="block",KERNEL=="sd*|mmcblk*" \
RUN+="/bin/sh -c '/sbin/make-fstab --update'"
